# jfa-go.service
[Unit]
Description=An account management system for Jellyfin.

[Service]
ExecStart={executable}

[Install]
WantedBy=default.target
