# data/services/jf-accounts.service
[Unit]
Description=A basic account management system for Jellyfin.

[Service]
ExecStart={executable}

[Install]
WantedBy=default.target
